# debian/gigi-proxy.socket
[Unit]
Description=Gigi Webserver module socket

[Socket]
ListenStream=80

[Install]
WantedBy=sockets.target
