# debian/gigi-proxy.service.d/Z-Disable-SystemCallFilter.conf
[Service]
# SystemCallFilter seems to be broken in some setups in Debian Stretch's systemd 232-r14
SystemCallFilter=
